# srv/munin/files/etc/cron.d/uws
*/1 * * * *	munin /usr/local/bin/munin-sendmail.py
