# srv/munin/files/etc/cron.d/uws
MAILTO=root
*/5 * * * *	munin /usr/local/bin/munin-limits.sh
