# srv/munin/files/etc/cron.d/uws
*/1 * * * *	munin /opt/munin/bin/munin-sendmail.py
