# srv/herokud/files/etc/cron.d/uws
USER=root
HOME=/root
MAILTO=root

@reboot	root echo 'herokud started'

# m h  dom mon dow	username command
*/2 * * * * uws /home/uws/etc/fetchmail.sh --syslog -K
