# srv/crond/files/etc/cron.d/uws
MAILTO=root

@reboot	root echo 'crond started'

# m h  dom mon dow	username command
