# host/assets/jsbatch/etc/cron.d/uwsbot
# cron-jobs for uwsbot
MAILTO=root

# bot/default api
1,6,11,16,21,26,31,36,41,46,51,56 * * * *	uwsrun cd /srv/uws/deploy && ./docker/uwsbot/run.sh -env bot/default -name api
